# deploy/fwlab-webui.service
[Unit]
Description=FW-LAB Web UI
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=cdomotor
Group=cdomotor
WorkingDirectory=/home/cdomotor/.openclaw/workspace/projects/ALERT1v3
ExecStart=/home/cdomotor/.openclaw/workspace/projects/ALERT1v3/tools/start_webui.sh
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
