# deploy/fwlab-log-retention.service
[Unit]
Description=FW-LAB log retention cleanup
After=network-online.target

[Service]
Type=oneshot
User=cdomotor
Group=cdomotor
WorkingDirectory=/home/cdomotor/.openclaw/workspace/projects/ALERT1v3
ExecStart=/usr/bin/python3 /home/cdomotor/.openclaw/workspace/projects/ALERT1v3/tools/log_retention.py --days 14 --max-mb 1024
